# etc/systemd/system/nscd-resolve-watch.path
[Unit]
Description=Watchdog for /etc/resolv.conf changes

[Path]
PathChanged=/etc/resolv.conf

[Install]
WantedBy=multi-user.target

